# deploy/vektra-trader.timer
[Unit]
Description=Kjor Vektra Trader hvert 15. minutt

[Timer]
OnCalendar=*:0/15
RandomizedDelaySec=60
Persistent=true

[Install]
WantedBy=timers.target
